# supervisord.conf
[supervisord]
nodaemon=true

[program:celery]
command=su -m project -c "celery -A app worker -l info"
autorestart=true